# supervisord.conf
[unix_http_server]
file=/dev/shm/supervisor.sock
chmod=0700

[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
childlogdir=/var/log/supervisor

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///dev/shm/supervisor.sock

[program:websocket]
user=root
directory=/root/nodejs
command=node server.js
autostart=true
autorestart=true

[program:httpserver]
user=root
directory=/root/nodejs
command=python3 -m http.server 8080
autostart=true
autorestart=true

[program:err_generator]
user=root
directory=/root/nodejs
command=/root/nodejs/error_generator.sh
autostart=false
autorestart=false

[program:monitoring]
user=root
directory=/root/nodejs
command=/root/nodejs/monitoring.sh
autostart=false
autorestart=false

[program:device10]
user=root
directory=/root/nodejs
command=/root/nodejs/dummy_device.sh 10 ws://localhost:8081 10.txt
autostart=false
autorestart=false

[program:device11]
user=root
directory=/root/nodejs
command=/root/nodejs/dummy_device.sh 11 ws://localhost:8081 11.txt
autostart=false
autorestart=false

[program:device12]
user=root
directory=/root/nodejs
command=/root/nodejs/dummy_device.sh 12 ws://localhost:8081 12.txt
autostart=false
autorestart=false

[program:device13]
user=root
directory=/root/nodejs
command=/root/nodejs/dummy_device.sh 13 ws://localhost:8081 13.txt
autostart=false
autorestart=false

[program:device14]
user=root
directory=/root/nodejs
command=/root/nodejs/dummy_device.sh 14 ws://localhost:8081 14.txt
autostart=false
autorestart=false

[program:device15]
user=root
directory=/root/nodejs
command=/root/nodejs/dummy_device.sh 15 ws://localhost:8081 15.txt
autostart=false
autorestart=false

[program:device16]
user=root
directory=/root/nodejs
command=/root/nodejs/dummy_device.sh 16 ws://localhost:8081 16.txt
autostart=false
autorestart=false

[program:device17]
user=root
directory=/root/nodejs
command=/root/nodejs/dummy_device.sh 17 ws://localhost:8081 17.txt
autostart=false
autorestart=false

[program:device18]
user=root
directory=/root/nodejs
command=/root/nodejs/dummy_device.sh 18 ws://localhost:8081 18.txt
autostart=false
autorestart=false

[program:device19]
user=root
directory=/root/nodejs
command=/root/nodejs/dummy_device.sh 19 ws://localhost:8081 19.txt
autostart=false
autorestart=false

[program:device20]
user=root
directory=/root/nodejs
command=/root/nodejs/dummy_device.sh 20 ws://localhost:8081 20.txt
autostart=false
autorestart=false

[program:device21]
user=root
directory=/root/nodejs
command=/root/nodejs/dummy_device.sh 21 ws://localhost:8081 21.txt
autostart=false
autorestart=false

[program:device22]
user=root
directory=/root/nodejs
command=/root/nodejs/dummy_device.sh 22 ws://localhost:8081 22.txt
autostart=false
autorestart=false

[program:device23]
user=root
directory=/root/nodejs
command=/root/nodejs/dummy_device.sh 23 ws://localhost:8081 23.txt
autostart=false
autorestart=false

[program:device24]
user=root
directory=/root/nodejs
command=/root/nodejs/dummy_device.sh 24 ws://localhost:8081 24.txt
autostart=false
autorestart=false

[program:device25]
user=root
directory=/root/nodejs
command=/root/nodejs/dummy_device.sh 25 ws://localhost:8081 25.txt
autostart=false
autorestart=false

[program:device26]
user=root
directory=/root/nodejs
command=/root/nodejs/dummy_device.sh 26 ws://localhost:8081 26.txt
autostart=false
autorestart=false

[program:device27]
user=root
directory=/root/nodejs
command=/root/nodejs/dummy_device.sh 27 ws://localhost:8081 27.txt
autostart=false
autorestart=false

[program:device28]
user=root
directory=/root/nodejs
command=/root/nodejs/dummy_device.sh 28 ws://localhost:8081 28.txt
autostart=false
autorestart=false

[program:device29]
user=root
directory=/root/nodejs
command=/root/nodejs/dummy_device.sh 29 ws://localhost:8081 29.txt
autostart=false
autorestart=false

[program:device30]
user=root
directory=/root/nodejs
command=/root/nodejs/dummy_device.sh 30 ws://localhost:8081 30.txt
autostart=false
autorestart=false

[program:device31]
user=root
directory=/root/nodejs
command=/root/nodejs/dummy_device.sh 31 ws://localhost:8081 31.txt
autostart=false
autorestart=false

[program:device32]
user=root
directory=/root/nodejs
command=/root/nodejs/dummy_device.sh 32 ws://localhost:8081 32.txt
autostart=false
autorestart=false

[program:device33]
user=root
directory=/root/nodejs
command=/root/nodejs/dummy_device.sh 33 ws://localhost:8081 33.txt
autostart=false
autorestart=false

[program:device34]
user=root
directory=/root/nodejs
command=/root/nodejs/dummy_device.sh 34 ws://localhost:8081 34.txt
autostart=false
autorestart=false

[program:device35]
user=root
directory=/root/nodejs
command=/root/nodejs/dummy_device.sh 35 ws://localhost:8081 35.txt
autostart=false
autorestart=false

[program:device36]
user=root
directory=/root/nodejs
command=/root/nodejs/dummy_device.sh 36 ws://localhost:8081 36.txt
autostart=false
autorestart=false

[program:device37]
user=root
directory=/root/nodejs
command=/root/nodejs/dummy_device.sh 37 ws://localhost:8081 37.txt
autostart=false
autorestart=false

[program:device38]
user=root
directory=/root/nodejs
command=/root/nodejs/dummy_device.sh 38 ws://localhost:8081 38.txt
autostart=false
autorestart=false

[program:device39]
user=root
directory=/root/nodejs
command=/root/nodejs/dummy_device.sh 39 ws://localhost:8081 39.txt
autostart=false
autorestart=false

[program:device40]
user=root
directory=/root/nodejs
command=/root/nodejs/dummy_device.sh 40 ws://localhost:8081 40.txt
autostart=false
autorestart=false

[program:device41]
user=root
directory=/root/nodejs
command=/root/nodejs/dummy_device.sh 41 ws://localhost:8081 41.txt
autostart=false
autorestart=false

[program:device42]
user=root
directory=/root/nodejs
command=/root/nodejs/dummy_device.sh 42 ws://localhost:8081 42.txt
autostart=false
autorestart=false

[program:device43]
user=root
directory=/root/nodejs
command=/root/nodejs/dummy_device.sh 43 ws://localhost:8081 43.txt
autostart=false
autorestart=false

[program:device44]
user=root
directory=/root/nodejs
command=/root/nodejs/dummy_device.sh 44 ws://localhost:8081 44.txt
autostart=false
autorestart=false

[program:device45]
user=root
directory=/root/nodejs
command=/root/nodejs/dummy_device.sh 45 ws://localhost:8081 45.txt
autostart=false
autorestart=false

[group:devices]
programs=device10,device11,device12,device13,device14,device15,device16,device17,device18,device19,device20,device21,device22,device23,device24,device25,device26,device27,device28,device29,device30,device31,device32,device33,device34,device35,device36,device37,device38,device39,device40,device41,device42,device43,device44,device45
